# debian/python-thumbor-wikimedia.postinst
#!/bin/sh

set -e

if [ "$1" = "configure" ] ; then
    invoke-rc.d thumbor restart || true
fi

#DEBHELPER#

exit 0
